# inc/supervisord.ini
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
autostart=true
autorestart=true

[program:mysql]
user=mysql
command=/usr/bin/pidproxy /var/run/mariadb/mariadb.pid /usr/bin/mysqld_safe
autostart=true
autorestart=true
